# config/systemd/geyser-webhook.service
[Unit]
Description=Geyser Webhook Service
After=geyser.service
Requires=geyser.service

[Service]
Type=simple
User=arkandias
WorkingDirectory=/home/arkandias/geyser-monorepo
ExecStart=/home/arkandias/geyser-monorepo/scripts/geyser webhook-start
ExecStopPost=/home/arkandias/geyser-monorepo/scripts/geyser webhook-stop
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target